# deploy/server/hydeclaw-core.service
[Unit]
Description=HydeClaw Core
After=network.target docker.service

[Service]
Type=simple
WorkingDirectory=/home/aronmav/hydeclaw
ExecStart=/home/aronmav/hydeclaw/hydeclaw-core-x86_64
EnvironmentFile=/home/aronmav/hydeclaw/.env
Environment=PATH=/home/aronmav/.bun/bin:/home/aronmav/.local/bin:/usr/local/bin:/usr/bin:/bin
# Egress via the xray HTTP proxy — this server's direct internet is region-blocked
# (Telegram etc. unreachable). Internal targets (localhost services, the openedai
# TTS at 10.10.1.42) bypass the proxy. Inherited by the toolgate/channels children.
Environment=HTTP_PROXY=http://10.10.1.37:1081
Environment=HTTPS_PROXY=http://10.10.1.37:1081
Environment=http_proxy=http://10.10.1.37:1081
Environment=https_proxy=http://10.10.1.37:1081
Environment=NO_PROXY=localhost,127.0.0.1,::1,10.10.1.42,deepseek.com,ollama.com,openrouter.ai,together.ai,fal.run,fal.ai
Environment=no_proxy=localhost,127.0.0.1,::1,10.10.1.42,deepseek.com,ollama.com,openrouter.ai,together.ai,fal.run,fal.ai
StandardOutput=append:/home/aronmav/hydeclaw/logs/core.log
StandardError=append:/home/aronmav/hydeclaw/logs/core.log
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
