# docker/opex-core.service
[Unit]
Description=Opex Core — AI Gateway
After=network.target docker.service
Wants=docker.service

[Service]
Type=notify
WatchdogSec=120
ExecStart=%h/opex/opex-core
WorkingDirectory=%h/opex
EnvironmentFile=%h/opex/.env
TimeoutStopSec=40
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryMax=256M
MemorySwapMax=512M

[Install]
WantedBy=multi-user.target
